# shell/profile.d/cvs.sh
# CVS Repository

CVSROOT=":pserver:${USER}@cvshost:/home/cvs/develop"
export CVSROOT
