# shell/profile.d/editor.sh
export EDITOR=/usr/bin/vi
export TEXEDIT='/usr/bin/vi %s'
