# shell/profile.d/editor.sh
export EDITOR=/usr/bin/vi
export PAGER=/usr/bin/less	
export TEXEDIT='/usr/bin/vi %s'
