# shell/profile.d/less.sh
## less ¤Î»ØÄê

LESSCHARSET=japanese-euc
export LESSCHARSET
