# shell/profile.d/postgresql.sh
## use PostgresSQL
#if [ -f /etc/postgresql/postgresql.env ]; then
#    . /etc/postgresql/postgresql.env
#fi
#
