# shell/profile.d/x.sh
# X Application Resource Directory
export XAPPLRESDIR=${CONFDIR}/X11/app-resources/
